# archlinux/kadnode.install

post_install() {
	printf "To enable service KadNode: systemctl enable kadnode\n"

	#add kadnode entry to NSS config file
	sed -i -e '/kadnode/!s/^\(hosts:.*\)dns\(.*\)/\1kadnode dns\2/' ${DESTDIR}/etc/nsswitch.conf
}

post_remove() {
	#remove kadnode entry from NSS config file
	sed -i -e 's/^\(hosts:.*\)kadnode \(.*\)/\1\2/' ${DESTDIR}/etc/nsswitch.conf
}
