# archlinux/kadnode.service
[Unit]
Description=P2P DNS service
After=network.target

[Service]
ExecStart=/usr/bin/kadnode --daemon --peerfile /etc/kadnode/peers.txt --config /etc/kadnode/kadnode.conf
Type=forking

[Install]
WantedBy=multi-user.target
